# macos/com.agentstack.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentstack</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>AGENTSTACK_DIR/start.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>AGENTSTACK_DIR</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>HOME_DIR/Library/Logs/agentstack/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>HOME_DIR/Library/Logs/agentstack/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
